# departments/platform/deployment/systemd/oqp-paper-snapshot.service
[Unit]
Description=Alpha Factory Paper Trading Snapshot Job
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/oqp_new
EnvironmentFile=-/home/ubuntu/.oqp_server_env
Environment=PYTHONPATH=src:.
ExecStart=/home/ubuntu/oqp_new/scripts/run_paper_snapshot_job.sh
StandardOutput=append:/home/ubuntu/oqp_new/logs/paper_snapshot_job.log
StandardError=append:/home/ubuntu/oqp_new/logs/paper_snapshot_job.log
