# departments/platform/deployment/systemd/oqp-paper-strategy-runner.service
[Unit]
Description=Alpha Factory Paper Strategy Runner Job
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/oqp_new
EnvironmentFile=-/home/ubuntu/.oqp_server_env
ExecStart=/home/ubuntu/oqp_new/scripts/trading/run_paper_strategy_runner_job.sh
StandardOutput=append:/home/ubuntu/oqp_new/runtime/logs/paper_strategy_runner.log
StandardError=append:/home/ubuntu/oqp_new/runtime/logs/paper_strategy_runner.log
